# docker/build/etc/supervisor.d/worker.conf
[group:akeneo]
programs=daemon
priority=25

[program:daemon]
command = /app/bin/console akeneo:batch:job-queue-consumer-daemon
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
